# server/configs/systemd/security-containers.service.in
[Unit]
Description=Security Containers Server
After=libvirtd.service
Requires=libvirtd.service
ConditionVirtualization=no

[Service]
Type=simple
ExecStart=${CMAKE_INSTALL_PREFIX}/bin/security-containers-server
Restart=on-failure
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
