# server/configs/systemd/vasum.socket.in
[Socket]
ListenStream=${VSM_SERVER_IPC_SOCKET_PATH}
SocketMode=0755
SmackLabelIPIn=*
SmackLabelIPOut=@
Service=vasum.service

[Install]
WantedBy=sockets.target
