# server/configs/systemd/vasum.service.in
[Unit]
Description=Vasum Server
ConditionVirtualization=no

[Service]
Type=simple
ExecStart=${CMAKE_INSTALL_PREFIX}/bin/vasum-server -r
Restart=on-failure
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
