# server/configs/systemd/vasum.socket
[Socket]
ListenStream=/run/vasum.socket
SocketMode=0755
SmackLabelIPIn=*
SmackLabelIPOut=@
Service=vasum.service

[Install]
WantedBy=sockets.target
